# deploy/aicodereview.service
[Unit]
Description=AI Code Review Assistant (FastAPI / Uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ai-code-review/server
EnvironmentFile=/opt/ai-code-review/server/.env
ExecStart=/opt/ai-code-review/server/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 5000
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
